# crontab
* * * * * cd && cd ../.. && cd dbt && sh run.sh
